# tools/email/mail-sync.service
[Unit]
Description=Pull new mail from Gmail via lieer

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'cd ~/.mail/gmail && gmi pull && notmuch new'
